# etc/udev/70-snap.gnome-system-monitor.rules
# This file is automatically generated.
# opengl
KERNEL=="nvhost-*", TAG+="snap_gnome-system-monitor_gnome-system-monitor"
# opengl
KERNEL=="nvmap", TAG+="snap_gnome-system-monitor_gnome-system-monitor"
# opengl
KERNEL=="renderD[0-9]*", TAG+="snap_gnome-system-monitor_gnome-system-monitor"
# opengl
KERNEL=="tegra_dc_[0-9]*", TAG+="snap_gnome-system-monitor_gnome-system-monitor"
# opengl
KERNEL=="tegra_dc_ctrl", TAG+="snap_gnome-system-monitor_gnome-system-monitor"
# opengl
KERNEL=="vchiq", TAG+="snap_gnome-system-monitor_gnome-system-monitor"
# opengl
KERNEL=="vcsm-cma", TAG+="snap_gnome-system-monitor_gnome-system-monitor"
# opengl
SUBSYSTEM=="drm", KERNEL=="card[0-9]*", TAG+="snap_gnome-system-monitor_gnome-system-monitor"
TAG=="snap_gnome-system-monitor_gnome-system-monitor", RUN+="/usr/lib/snapd/snap-device-helper $env{ACTION} snap_gnome-system-monitor_gnome-system-monitor $devpath $major:$minor"
