# etc/udev/70-snap.gnome-calculator.rules
# This file is automatically generated.
# opengl
KERNEL=="nvhost-*", TAG+="snap_gnome-calculator_gnome-calculator"
# opengl
KERNEL=="nvmap", TAG+="snap_gnome-calculator_gnome-calculator"
# opengl
KERNEL=="renderD[0-9]*", TAG+="snap_gnome-calculator_gnome-calculator"
# opengl
KERNEL=="tegra_dc_[0-9]*", TAG+="snap_gnome-calculator_gnome-calculator"
# opengl
KERNEL=="tegra_dc_ctrl", TAG+="snap_gnome-calculator_gnome-calculator"
# opengl
KERNEL=="vchiq", TAG+="snap_gnome-calculator_gnome-calculator"
# opengl
KERNEL=="vcsm-cma", TAG+="snap_gnome-calculator_gnome-calculator"
# opengl
SUBSYSTEM=="drm", KERNEL=="card[0-9]*", TAG+="snap_gnome-calculator_gnome-calculator"
TAG=="snap_gnome-calculator_gnome-calculator", RUN+="/usr/lib/snapd/snap-device-helper $env{ACTION} snap_gnome-calculator_gnome-calculator $devpath $major:$minor"
