# etc/udev/70-snap.gnome-logs.rules
# This file is automatically generated.
# opengl
KERNEL=="nvhost-*", TAG+="snap_gnome-logs_gnome-logs"
# opengl
KERNEL=="nvmap", TAG+="snap_gnome-logs_gnome-logs"
# opengl
KERNEL=="renderD[0-9]*", TAG+="snap_gnome-logs_gnome-logs"
# opengl
KERNEL=="tegra_dc_[0-9]*", TAG+="snap_gnome-logs_gnome-logs"
# opengl
KERNEL=="tegra_dc_ctrl", TAG+="snap_gnome-logs_gnome-logs"
# opengl
KERNEL=="vchiq", TAG+="snap_gnome-logs_gnome-logs"
# opengl
KERNEL=="vcsm-cma", TAG+="snap_gnome-logs_gnome-logs"
# opengl
SUBSYSTEM=="drm", KERNEL=="card[0-9]*", TAG+="snap_gnome-logs_gnome-logs"
TAG=="snap_gnome-logs_gnome-logs", RUN+="/usr/lib/snapd/snap-device-helper $env{ACTION} snap_gnome-logs_gnome-logs $devpath $major:$minor"
