# etc/udev/95-pty.rules
ACTION=="add" \
, KERNEL=="sd*" \
, SUBSYSTEM=="block" \
, KERNELS=="4-10.4.3.1" \
, SYMLINK+="sdcard0" \
, RUN+="/bin/chmod 777 /dev/$name"
